# smtp-tunnel.service
[Unit]
Description=SMTP 隧道服务器
Documentation=https://github.com/user/smtp-tunnel
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/smtp-tunnel
ExecStart=/usr/bin/python3 /opt/smtp-tunnel/server.py -c /etc/smtp-tunnel/config.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# 安全加固（可选）
# NoNewPrivileges=true      # 禁止获取新权限
# ProtectSystem=strict      # 严格保护系统目录
# ProtectHome=true          # 保护用户主目录
# ReadWritePaths=/etc/smtp-tunnel  # 允许写入的路径

[Install]
WantedBy=multi-user.target
